# smokeping_prober.service.in
[Unit]
Description=Export smokeping to Prometheus

[Service]
EnvironmentFile=-@SYSCONFDIR@/default/@NAME@
ExecStart=@BINDIR@/@NAME@ $ARGS

[Install]
WantedBy=multi-user.target
